# etc/systemd/system/kvm-server.socket
[Unit]
Description=KVM Server socket

[Socket]
ListenStream=/var/run/kvm-server.sock

[Install]
WantedBy=default.target
