# etc/systemd/system/kvm-server.service
[Unit]
Description=KVM Server

[Service]
ExecStart=/usr/local/bin/kvm-server

[Install]
WantedBy=multi-user.target
